# kubeconfig.path
[Path]
PathChanged = %h/.kube/configs

[Install]
WantedBy=default.target
