# kubeconfig.service
[Service]
Type=oneshot
ExecStart = /usr/bin/make -C %h/.kube
